# package/hpd.service
[Unit]
Description=Handheld Power Daemon (hpd)
Documentation=https://github.com/CiroDev-Git/hpd-handheld-power-daemon
After=multi-user.target dbus.service power-profiles-daemon.service
Wants=dbus.service
# hpd is the platform's single power authority. power-profiles-daemon
# writes the same platform_profile/EPP surfaces, so refuse to co-run.
# steamos-manager is D-Bus-activated and handled by `hpdctl doctor --fix`.
#
# The AUR post_install hook masks PPD (blocks both systemd and D-Bus
# activation). After= + Conflicts= remain here as a safety net for
# manual installs (install.sh) that do not run the pacman hooks.
#
# DELIBERATELY scoped to PPD only — do NOT add tuned or hhd here, even
# though hpd detects them as hard rivals (hpd_dbus::conflicts):
#   * Conflicts= is symmetric, so it only helps against a rival that is
#     ALSO masked. tuned is D-Bus-activatable (Type=dbus,
#     BusName=com.redhat.tuned): an unmasked `Conflicts=tuned.service`
#     lets the bus revive tuned and kill *hpd* — the exact v2.2.2
#     regression. With a mask it is redundant.
#   * hhd cannot be targeted at all: it runs as the templated
#     hhd@<user>.service whose instance name is unknown at package time.
#   * Unlike PPD (a headless service safe to mask silently), tuned and hhd
#     are user-chosen stacks. Neutralizing them is opt-in and informed via
#     `hpdctl doctor --fix`, never an automatic unit Conflicts= or hook.
Conflicts=power-profiles-daemon.service

[Service]
Type=simple
ExecStart=/usr/local/bin/hpd-daemon
# `systemctl reload hpd` triggers SIGHUP, which the daemon catches and
# uses to reload /etc/hpd/config.toml in-place (no process restart).
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
Environment="RUST_LOG=hpd=info,warn"

# Persistent state directory: systemd creates /var/lib/hpd with the right
# permissions and exposes its path through the STATE_DIRECTORY env var.
StateDirectory=hpd
StateDirectoryMode=0700

# Configuration directory: systemd ensures /etc/hpd exists and exposes
# its path via CONFIGURATION_DIRECTORY. The daemon reads
# `<dir>/config.toml`; missing or corrupt → falls back to defaults.
ConfigurationDirectory=hpd
ConfigurationDirectoryMode=0755

# --- Sandboxing ---
# hpd is root-only and writes to a small, well-known set of sysfs paths.
# Lock down everything else so a future bug can't escape the daemon.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateNetwork=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectClock=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# ProtectSystem=strict makes / read-only; we explicitly allow the sysfs
# leaves the daemon actually writes to.
ReadWritePaths=/sys/class/firmware-attributes/ /sys/class/power_supply/ /sys/firmware/acpi/

[Install]
WantedBy=multi-user.target
